# config/systemd/user/hyprnightlight.service
[Unit]
Description=Hyprland Night Light

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/nightlight.sh
