# config/systemd/user/hyprnightlight.service
[Unit]
Description=Hyprland nightlight auto adjust

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/nightlight.sh